# samples/avorsmtp.supervisor.conf
[program:avorsmtp]
command=/opt/avorsmtp/current/avorsmtp
directory=/opt/avorsmtp/current/
user=root
autostart=true
autorestart=true