# samples/avorsmtp.supervisor.conf
[program:vorimport]
command=/opt/avorf1com/current/avorf1com
direcory=/opt/avorf1com/current/
user=root
autostart=true
autorestart=true